# packages/pimelon/pimelon-6.16.1.tar.gz/pimelon-6.16.1/pine/config/templates/systemd/pimelon-redis-socketio.service
[Unit]
Description="{{ pine_name }}-redis-socketio"
PartOf={{ pine_name }}-redis.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ redis_server }} {{ redis_socketio_config }}
StandardOutput=file:{{ pine_dir }}/logs/redis-socketio.log
StandardError=file:{{ pine_dir }}/logs/redis-socketio.error.log
WorkingDirectory={{ sites_dir }}
