# packages/wiggler-studiorabota/wiggler_studiorabota-0.0.18.tar.gz/wiggler_studiorabota-0.0.18/wiggler/start_wiggler.service
[Unit]
Description=WiggleR service

[Service]
ExecStart=/bin/bash /usr/bin/start_wiggler.sh

[Install]
WantedBy=multi-user.target