# packages/es7s/es7s-3.48.0.dev0-py3-none-any.whl/es7s/data/systemd/es7s.service
[Unit]
Description=es7s/daemon
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
RuntimeDirectory=es7s
ExecStart=/home/%USER/.local/bin/es7s-daemon
User=%UID

[Install]
WantedBy=multi-user.target es7s-indicators.service
