# packages/pimelon/pimelon-6.16.1-py3-none-any.whl/pine/config/templates/systemd/pimelon-redis.target
[Unit]
After=network.target
Wants={{ pine_name }}-redis-cache.service {{ pine_name }}-redis-queue.service {{ pine_name }}-redis-socketio.service

[Install]
WantedBy=multi-user.target
