# packages/wiggler-studiorabota/wiggler_studiorabota-0.0.17.tar.gz/wiggler_studiorabota-0.0.17/wiggler/start_wiggler.service
[Unit]
Description=WiggleR service

[Service]
WorkingDirectory=/home/vincent/
User=vincent
ExecStart=/bin/bash /usr/bin/start_wiggler.sh

[Install]
WantedBy=multi-user.target